# config/paludis/local-repo/packages/mail-mta/postfix/files/systemd/postfix.service
[Unit]
Description=postfix
After=network.target

[Service]
Type=forking
ExecStart=/usr/host/bin/postfix start
ExecStop=/usr/host/bin/postfix stop
ExecReload=/usr/host/bin/postfix reload

[Install]
WantedBy=multi-user.target
